# roles/burp_client_configuration/templates/etc/cron.d/burp
#!/usr/bin/env bash
{% set cdata = cops_burpclient_vars %}
MAILTO="{{cdata.mailto}}"
{{cdata.crons}}
