# services_virt_docker/templates/etc/systemd/system/docker.service.d/ansible.conf
# {{ ansible_managed }}
[Service]
{% if corpusops_services_virt_docker_settings_http_proxy %}
Environment=HTTP_PROXY={{corpusops_services_virt_docker_settings_http_proxy}}
Environment=http_proxy={{corpusops_services_virt_docker_settings_http_proxy}}
{% endif %}
{% if corpusops_services_virt_docker_settings_TMPDIR %}
Environment=TMPDIR={{corpusops_services_virt_docker_settings_TMPDIR}}
{% endif %}
ExecStart=
ExecStart=/usr/bin/dockerd{% if corpusops_services_virt_docker_settings_DOCKERD_OPTS%} {{corpusops_services_virt_docker_settings_DOCKERD_OPTS}}{%endif%}
Type=notify
