# roles/services_virt_lxc/templates/etc/systemd/system/lxc.service.d/zz_cops_lxc.conf
[Unit]
After=syslog.target network.target lxc-net.service lxc-net-makina.service cops_magicbridge_copslxcbr.service
Wants=lxc-net.service cops_magicbridge_copslxcbr.service
