# roles/services_virt_docker/templates/etc/systemd/system/docker.service.d/ansible.conf
# {{ ansible_managed }}
[Service]
ExecStart=
ExecStart=/bin/sh -c '\
  found=;\
  for i in /etc/default /etc/sysconfig;do\
   for j in $i/docker $i/docker-storage $i/docker-network;do\
    if [ -e $j ];then . $j;found=1;fi;\
   done; \
   if [ "x$found" = "x1" ];then echo "Loaded from $i";break;fi;\
  done;\
  exec $DOCKERD \
    $DOCKER_OPTS \
    $DOCKER_STORAGE_OPTIONS \
    $DOCKER_NETWORK_OPTIONS \
    $BLOCK_REGISTRY $INSECURE_REGISTRY'
