# docker_compose_service/templates/upstart.conf
# {{ansible_managed}}
#{% set d = docker_compose_service_vars.project %}
#{% set w = docker_compose_service_vars.working_directory %}
description "{{d}} DockerCompose service"
author "corpusops"
start on filesystem and started docker
stop on runlevel [!2345]
# Automatically restart process if crashed (limit 3 times in 4 minutes)
respawn
kill signal SIGINT
script
set -x
cd "{{w}}"
/usr/bin/env docker-compose -p "{{d}}" pull
/usr/bin/env docker-compose -p "{{d}}" up
end script

post-stop script
set -x
cd "{{w}}"
/usr/bin/env docker-compose -p "{{d}}" stop || /bin/true
/usr/bin/env docker-compose -p "{{d}}" down || /bin/true
end script
# vim: set ft=upstart :
