# docker_compose_service/templates/upstart.conf
# {{ansible_managed}}

#{% set v = docker_compose_service_vars %}
#{% set d = v.project %}
#{% set w = v.path %}
description "{{d}} DockerCompose service"
author "corpusops"
start on filesystem and started docker
stop on runlevel [!2345]
# Automatically restart process if crashed (limit 3 times in 4 minutes)
{% if "{0}".format(v.restart).lower() not in ["no", "false"] %}
respawn
{% endif %}
kill signal SIGINT

script
set -x
cd "{{w}}"
/usr/bin/env docker-compose -p "{{d}}" {% for i in v._files%}-f "{{i}}" {%endfor%} {{v.pull_args}} pull &&\
/usr/bin/env docker-compose -p "{{d}}" {% for i in v._files%}-f "{{i}}" {%endfor%} {{v.up_args}} up
end script

pre-start script
cd "{{w}}"
/usr/bin/env docker-compose -p "{{d}}" {% for i in v._files%}-f "{{i}}" {%endfor%} {{v.config_args}} config || exit 1
end script

pre-stop script
cd "{{w}}"
/usr/bin/env docker-compose -p "{{d}}" {% for i in v._files%}-f "{{i}}" {%endfor%} {{v.config_args}} config || exit 1
end script

post-stop script
set -x
cd "{{w}}"
/usr/bin/env docker-compose -p "{{d}}" {% for i in v._files%}-f "{{i}}" {%endfor%} {{v.stop_args}} stop || /bin/true
/usr/bin/env docker-compose -p "{{d}}" {% for i in v._files%}-f "{{i}}" {%endfor%} {{v.down_args}} down || /bin/true
end script
# vim: set ft=upstart :
